# systemd/server-watchdog-monthly.service
[Unit]
Description=server-watchdog Monthly Maintenance Report
Documentation=https://github.com/matebuteler/server-watchdog
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/server-watchdog-monthly

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=server-watchdog-monthly

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
# Allow writes to the watchdog log dir and common package-manager log files.
# Non-existent paths are silently ignored by systemd, so these entries are
# safe on systems that use apt or zypper instead of dnf.
ReadWritePaths=/var/log/server-watchdog /var/log/dnf.log /var/log/dnf.librepo.log /var/log/dnf.rpm.log /var/log/apt
